# orion/systemd/orion-runner.service
[Unit]
Description=Mega Orion runner (Buck2 worker)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=orion
Group=orion
WorkingDirectory=/home/orion/orion-runner

# Optional runtime config. Create this file on the VM to set SERVER_WS,
# BUCK_PROJECT_ROOT, ORION_WORKER_ID, etc.
EnvironmentFile=-/etc/orion/orion-runner.env

# Scorpiofs config path (stable default expected by orion/src/antares.rs)
Environment=SCORPIO_CONFIG=/home/orion/orion-runner/scorpio.toml

# Cleanup before start: kill old processes and unmount FUSE
ExecStartPre=/home/orion/orion-runner/cleanup.sh

# Start orion (run.sh loads .env and starts the binary)
ExecStart=/home/orion/orion-runner/run.sh
Restart=always
RestartSec=2

# Grant capabilities for FUSE mounting
AmbientCapabilities=CAP_SYS_ADMIN
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_DAC_OVERRIDE

# Basic hardening (relaxed for FUSE operations)
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
